# revista-expats.service
[Unit]
Description=Revista Expats AI Server
After=network.target

[Service]
Type=simple
User=fleet
WorkingDirectory=/home/fleet/Escritorio/Carepetas proyects/Revista-expats-ai
ExecStart=/home/fleet/Escritorio/Carepetas proyects/Revista-expats-ai/.venv/bin/uvicorn main:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=5
Environment=PRODUCTION=true

[Install]
WantedBy=multi-user.target